# systemd/user/plasma-xdg-desktop-portal-kde.service.d/override.conf
[Service]
Environment="QT_QPA_PLATFORMTHEME=qt6ct"

